# deploy/systemd/llm-benchmark-ollama.service
[Unit]
Description=Isolated Ollama backend for local LLM benchmarks
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=llm-benchmark
Group=llm-benchmark
SupplementaryGroups=video render
WorkingDirectory=/var/lib/llm-benchmark
Environment=HOME=/var/lib/llm-benchmark
Environment=OLLAMA_HOST=127.0.0.1:11435
Environment=OLLAMA_MODELS=/var/lib/llm-benchmark/models
Environment=OLLAMA_KEEP_ALIVE=30s
Environment=OLLAMA_MAX_LOADED_MODELS=1
Environment=OLLAMA_NUM_PARALLEL=1
ExecStart=/usr/local/bin/ollama serve
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/llm-benchmark
RestrictSUIDSGID=true
LockPersonality=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
