# deploy/systemd/nova-voice-tts.service
[Unit]
Description=Nova Voice streaming Qwen3-TTS runtime
After=network-online.target nova-voice-llm.service
Wants=network-online.target
Requires=nova-voice-llm.service

[Service]
Type=simple
User=antidamage
Group=antidamage
WorkingDirectory=/home/antidamage/nova-tts-stream
Environment=HOME=/home/antidamage
Environment=HF_HOME=/home/antidamage/.cache/huggingface
Environment=HF_HUB_OFFLINE=1
Environment=TRANSFORMERS_OFFLINE=1
Environment=VLLM_USE_FLASHINFER_SAMPLER=0
Environment=CUDA_VISIBLE_DEVICES=0
Environment=PYTHONUNBUFFERED=1
# Warm the Qwen3-TTS weights into page cache ahead of the vLLM load so the
# runtime reads them from RAM. Best-effort; the '-' tolerates a missing binary.
ExecStartPre=-/usr/bin/vmtouch -qt /home/antidamage/.cache/huggingface/hub/models--Qwen--Qwen3-TTS-12Hz-0.6B-CustomVoice
ExecStart=/home/antidamage/nova-tts-stream/venv/bin/vllm-omni serve Qwen/Qwen3-TTS-12Hz-0.6B-CustomVoice --omni --deploy-config /home/antidamage/nova-tts-stream/qwen3_tts_turing.yaml --host 127.0.0.1 --port 8091 --trust-remote-code
ExecStartPost=/usr/bin/bash /opt/nova-voice/current/ops/warm-streaming-tts.sh
# Soft memory ceiling (see nova-voice-llm.service). vLLM's host-side RSS plus
# pinned staging fits well under this on the 24 GiB host; it only catches a
# runaway before the system-wide earlyoom backstop can fire.
MemoryAccounting=yes
MemoryHigh=8G
Restart=on-failure
RestartSec=5
TimeoutStartSec=180
TimeoutStopSec=45
KillMode=mixed
LimitMEMLOCK=infinity
TasksMax=infinity
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/antidamage/.cache
RestrictSUIDSGID=true
LockPersonality=true
# Gloo enumerates host interfaces through netlink before opening its local
# single-rank TCP process group.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
IPAddressDeny=any
IPAddressAllow=localhost
IPAddressAllow=192.168.8.0/24

[Install]
WantedBy=multi-user.target
